# deploy/systemd/creativesuite-frontend.service
[Unit]
Description=CreativeSuite ERP Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/creativesuite-frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start -- -p 3000 -H 127.0.0.1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target