# deploy/systemd/creativesuite-scheduler.service
[Unit]
Description=CreativeSuite ERP Scheduler
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/creativesuite-erp
ExecStart=/usr/bin/php artisan schedule:work
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target